# blueprints_ng/modules/sdcore_upf/config/router.service.jinja2
[Unit]
Description=UPF ROUTER {{ side }}
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/router.sh

[Install]
WantedBy=multi-user.target
