# src/nfvcl/blueprints_ng/modules/ueransim/config/config-network.service.jinja2
[Unit]
Description=Configuring GNB Network
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/config_network.sh

[Install]
WantedBy=multi-user.target
