# src/nfvcl/blueprints_ng/modules/oai/oai_upf/upf_forward.service.jinja2
[Unit]
Description=Enable forward for UPF
After=network.target
StartLimitIntervalSec=0

[Service]
Type=oneshot
User=root
ExecStart=/opt/upf_forward.sh

[Install]
WantedBy=multi-user.target
